# crontab
# Run logrotate.sh daily at midnight
0 * * * * /app/logrotate.sh >> /var/log/cron.log 2>&1
0 * * * * echo "Hello World! I am working!" >> /var/log/cron.log 2>&1
